# services/task-sync.timer
[Unit]
Description=Run taskwarrior sync every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target

