# services/task-sync.service
[Unit]
Description=Run taskwarrior sync

[Service]
Type=oneshot
ExecStart=/home/{{ username }}/.local/bin/task-sync.sh
